# server/proxyguard-frontend.service

[Unit]
Description=ProxyGuard Frontend Application
After=network.target proxyguard-api.service
Requires=proxyguard-api.service

[Service]
Type=simple
User=proxyguard
Group=proxyguard
WorkingDirectory=/opt/proxyguard
ExecStart=/usr/bin/node /opt/proxyguard/start-frontend.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=proxyguard-frontend
Environment=NODE_ENV=production PORT=3000

[Install]
WantedBy=multi-user.target
